# ansible/playbooks/router/firewall_script.service
[Unit]
Description=Firewall Script
After=network.target
[Service]
ExecStart=/root/firewall.sh
[Install]
WantedBy=default.target
